# roles/calc-sum/templates/goapp.service
[Unit]
Description=goapp

[Service]
Type=simple
ExecStart=/srv/main
ExecStop=/bin/kill $MAINPID
